# misc/emeter.service
[Unit]
Description=Emeter
After=network.target

[Service]
User=pi
ExecStart=python /home/pi/reefwatcher/emeter.py
