# misc/reef.service
[Unit]
Description=Reef Watcher
After=network.target

[Service]
User=pi
Restart=always
RestartSec=15s
ExecStart=python /home/pi/reefwatcher/reef.py
